# devops/ansible/templates/supervisor.ini
[program:{{ project_name }}]
command={{ virtual_env }}/bin/uwsgi
  --socket /tmp/{{ project_name }}.sock
  --logto {{ projects_root }}/log/{{ project_name }}/application.log
  --home {{ virtual_env }}
  --pythonpath {{ project_folder }}/server/
  --wsgi-file {{ project_folder }}/server/main.py
  --callable app
  --max-requests 1000
  --master
  --processes 1
  --chmod

directory={{ project_folder }}/server
autostart=true
autorestart=true
